# conf/supervisord.d/GetIdGTmetrix.ini
[program:GetIdGTmetrix]
command=/usr/bin/php /srv/www/bin/cli.php v1 Worker listen GetIdGTmetrix
autostart=true
autorestart=true
numprocs=1
startsecs=5
environment=VERBOSE="1"
stderr_logfile=/var/log/GetIdGTmetrix-err.log
stdout_logfile=/var/log/GetIdGTmetrix-out.log
user=root
