# conf/supervisord.d/GetPerformanceGTmetrix.ini
[program:GetPerformanceGTmetrix]
command=/usr/bin/php /srv/www/bin/cli.php v1 Worker listen GetPerformanceGTmetrix
autostart=true
autorestart=true
numprocs=1
startsecs=5
environment=VERBOSE="1"
stderr_logfile=/var/log/GetPerformanceGTmetrix-err.log
stdout_logfile=/var/log/GetPerformanceGTmetrix-out.log
user=root
