# conf/supervisord.d/SendingEmail.ini
[program:SendingEmail]
command=/usr/bin/php /srv/www/bin/cli.php v1 Worker listen SendingEmail
autostart=true
autorestart=true
numprocs=1
startsecs=5
environment=VERBOSE="1"
stderr_logfile=/var/log/SendingEmail-err.log
stdout_logfile=/var/log/SendingEmail-out.log
user=root
